# helix-swarm/etc/supervisor/conf.d/cron.conf
[program:cron]
autorestart=false
command=cron -f
