# helix-swarm/etc/supervisor/conf.d/helix-swarm.conf
[program:helix-swarm]
command=/usr/sbin/apachectl -D FOREGROUND
