# helix-swarm/etc/supervisor/conf.d/redis-server.conf
[program:redis-server]
command=/usr/bin/redis-server /opt/perforce/etc/redis-server.conf
